# etc/profile.d/01-java.sh
JAVA_HOME=/java/jdk
GROOVY_HOME=/java/groovy
ANT_HOME=/java/ant

PATH=${JAVA_HOME}/bin:${PATH}:${GROOVY_HOME}/bin:${ANT_HOME}/bin:/java/gradle/bin:/java/maven/bin
